# root/etc/services.d/postfix/run
#!/usr/bin/with-contenv sh


exec postfix start
