# root/etc/periodic/15min/01-drupal
#!/bin/bash

/usr/bin/env PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin COLUMNS=72 /usr/local/drush/drush --root=/DATA/htdocs --uri=$VIRTUAL_HOST --quiet cron
